# bibleai-backoffice.service
[Unit]
Description=BibleAI Backoffice Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/bibleai
ExecStart=/home/ec2-user/bibleai/backoffice
Restart=always
RestartSec=5
StandardOutput=append:/home/ec2-user/bibleai/backoffice.log
StandardError=append:/home/ec2-user/bibleai/backoffice.log

# 환경 변수
Environment="PORT=9090"
Environment="TZ=Asia/Seoul"

[Install]
WantedBy=multi-user.target
